# docker-sickrage/root/etc/services.d/sickrage/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc python3 /app/sickrage/SickChill.py \
	--datadir=/config --port 8081
